# debian/perspective-digital-ocean.upstart
#!/usr/bin/env bash

description "perspective-digital-ocean"

start on runlevel [2345]
stop on runlevel [06]

limit nofile 4096 4096
setuid perspective

exec /usr/share/perspective/perspective-digital-ocean/start.sh
