# debian/perspective-storage.upstart
#!/usr/bin/env bash

description "perspective-storage"

start on runlevel [2345]
stop on runlevel [06]

limit nofile 4096 4096
setuid perspective

exec /usr/share/perspective/perspective-storage/start.sh
